# root/etc/services.d/php-fpm/run
#!/bin/execlineb -P
fdmove -c 2 1
exec -c
s6-envuidgid www-data /usr/local/sbin/php-fpm --nodaemonize
